# src-tauri/helpers/com.alarmmaster.wake-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alarmmaster.wake-helper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/AlarmMasterWakeHelper</string>
    </array>

    <!-- Trigger when schedule file changes -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/Shared/AlarmMaster/schedule.json</string>
    </array>

    <!-- Also run once per hour as safety net -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Don't run at boot, only when triggered -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/Shared/AlarmMaster/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/AlarmMaster/launchd-stderr.log</string>
</dict>
</plist>
